# config/mnt-data.automount
[Unit]
Description="automount /dev/sda1 /mnt/data"

[Automount]
Where=/mnt/data
TimeoutIdleSec=1h 30m

[Install]
WantedBy=multi-user.target
